# templates/etc/profile.d/safe-rm.sh.erb
# MANAGED BY PUPPET
# creates an alias for the rm command
alias rm="<%= @bin_path %>safe-rm"
